# test/inittab
::sysinit:/etc/init.d/rcS
::askfirst:-/bin/sh
